# nginx-php56/supervisord.conf
[program:php5-fpm]
command = php5-fpm --fpm-config /etc/php5/fpm/php-fpm.conf
stdout_logfile = /data/logs/php-fpm.log
stderr_logfile = /data/logs/php-fpm.log

