# php72-nginx-alpine/etc/before_launch.ini
[program:before_launch]
command:sh /before-launch.sh
autorestart = false
startretries = 1
startsecs = 0
