# nginx-php/supervisord.conf
[program:php-fpm7.1]
command = php-fpm7.1 --fpm-config /etc/php/7.1/fpm/php-fpm.conf
stdout_logfile = /data/logs/php-fpm.log
stderr_logfile = /data/logs/php-fpm.log

